# launchd/com.cn.feishu-bot-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cn.feishu-bot-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-u</string>
    <string>/Users/cn/Workspace/feishu-bot-bridge/ws_bot.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/cn/Workspace/feishu-bot-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PATH</key>
    <string>/Users/cn/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/cn</string>
    <key>CLAUDE_CODE_OAUTH_TOKEN</key>
    <string>YOUR_CLAUDE_CODE_OAUTH_TOKEN_HERE</string>
    <key>http_proxy</key>
    <string>http://127.0.0.1:1089</string>
    <key>https_proxy</key>
    <string>http://127.0.0.1:1089</string>
    <key>all_proxy</key>
    <string>socks5://127.0.0.1:1084</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/cn/Workspace/feishu-bot-bridge/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/cn/Workspace/feishu-bot-bridge/logs/launchd.err.log</string>
</dict>
</plist>
