# launchd/com.cn.feishu-daily-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cn.feishu-daily-report</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-u</string>
    <string>/Users/cn/Workspace/feishu-bot-bridge/scripts/daily_report_job.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/cn/Workspace/feishu-bot-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/cn/Workspace/feishu-bot-bridge/logs/daily-report.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/cn/Workspace/feishu-bot-bridge/logs/daily-report.err.log</string>
</dict>
</plist>
